# systemd/restic-backup@.service
[Unit]
Description=Run Restic backup job %i
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/docker-restic-config/systemd.env
ExecStart=/bin/sh -eu -c 'cd "$${STACK_DIR}" && exec docker compose --project-directory "$${STACK_DIR}" -f "$${COMPOSE_FILE}" run --rm restic-job /scripts/restic-job.sh "%i"'
StandardOutput=journal
StandardError=journal
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
